# ops/signage-overlay.service
[Unit]
Description=ubuntuapp Digital Signage Overlay
Documentation=file:///opt/ubuntuapp/docs/DEPLOY-AND-VERIFY.md
After=graphical-session.target
Wants=graphical-session.target

# 60 秒以内に 5 回を超えてクラッシュした場合は起動を停止する
# 復旧: systemctl --user reset-failed signage-overlay && systemctl --user start signage-overlay
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# 表示サーバ環境変数を明示設定
# X11 環境では DISPLAY=:0、Wayland 環境では WAYLAND_DISPLAY=wayland-0 が有効になる
# XDG_RUNTIME_DIR は UID 1000 の場合の既定値。別ユーザーの場合は install.sh が書き換える
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000

# AppImage のフルパス。install.sh が実際のパスに書き換える
ExecStart=/opt/ubuntuapp/ubuntuapp.AppImage

# 異常終了後 5 秒で自動再起動
Restart=always
RestartSec=5

# 標準出力・エラー出力を journald へ転送
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signage-overlay

[Install]
WantedBy=graphical-session.target
